# ansible/devOps-project/certbot/files/certbot.service
[Unit]
Description=Lets Encrypt renewal

[Service]
Type=oneshot
ExecStart=/root/certbot/get-certs.sh
